# DEBIAN/postinst
systemctl daemon-reload
systemctl enable jxcore.service